# files/common/etc/profile.d/editor.sh
# Ensure vim is the default editor for interactive shells
export EDITOR=vim
export VISUAL=vim
