# files/personal/etc/systemd/system/sshd.service.d/disable.conf
[Unit]
ConditionPathExists=!/etc/ssh/enable-sshd
